# LEMP/ubuntu/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/sbin/php-fpm7.2 -F
autostart=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
